# codis-image/conf/supervisord_slave.conf
[supervisord]
nodaemon=true

[program:slave]
command=/bin/bash /opt/codis/bin/codis-start rediss start
autorestart=false
redirect_stderr=true
stdout_logfile=/opt/codis/codis-start.log
stopasgroup=true
